# sicora-be-go/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/sicora/supervisord.log
pidfile=/var/run/supervisord.pid

[program:userservice]
command=/app/userservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/userservice.log
stderr_logfile=/var/log/sicora/userservice_error.log
environment=PORT="8080",SERVICE_NAME="userservice"

[program:attendanceservice]
command=/app/attendanceservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/attendanceservice.log
stderr_logfile=/var/log/sicora/attendanceservice_error.log
environment=PORT="8081",SERVICE_NAME="attendanceservice"

[program:scheduleservice]
command=/app/scheduleservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/scheduleservice.log
stderr_logfile=/var/log/sicora/scheduleservice_error.log
environment=PORT="8082",SERVICE_NAME="scheduleservice"

[program:evalinservice]
command=/app/evalinservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/evalinservice.log
stderr_logfile=/var/log/sicora/evalinservice_error.log
environment=PORT="8083",SERVICE_NAME="evalinservice"

[program:mevalservice]
command=/app/mevalservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/mevalservice.log
stderr_logfile=/var/log/sicora/mevalservice_error.log
environment=PORT="8084",SERVICE_NAME="mevalservice"

[program:projectevalservice]
command=/app/projectevalservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/projectevalservice.log
stderr_logfile=/var/log/sicora/projectevalservice_error.log
environment=PORT="8085",SERVICE_NAME="projectevalservice"

[program:kbservice]
command=/app/kbservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/kbservice.log
stderr_logfile=/var/log/sicora/kbservice_error.log
environment=PORT="8086",SERVICE_NAME="kbservice"
