# sicora-be-go/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/sicora/supervisord.log
pidfile=/var/run/supervisord.pid

[program:attendanceservice]
command=/app/attendanceservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/attendanceservice.log
stderr_logfile=/var/log/sicora/attendanceservice_error.log
environment=PORT="8080",SERVICE_NAME="attendanceservice"

[program:projectevalservice]
command=/app/projectevalservice
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sicora/projectevalservice.log
stderr_logfile=/var/log/sicora/projectevalservice_error.log
environment=PORT="8081",SERVICE_NAME="projectevalservice"

# TODO: Agregar más servicios cuando estén listos:
# - userservice: requiere swag docs
# - scheduleservice: archivos vacíos
# - evalinservice: sin docs/
# - mevalservice: sin docs/
# - kbservice: FAQStats undefined
